# installers/com.ttyd.terminal.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ttyd.terminal</string>
    <key>ProgramArguments</key>
    <array>
        <string>__TTYD_BIN__</string>
        <string>-W</string>
        <string>-p</string>
        <string>7681</string>
        <string>-i</string>
        <string>127.0.0.1</string>
        <string>login</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/log/ttyd.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/ttyd.log</string>
</dict>
</plist>
